# systemd/vibmon-gateway.service
[Unit]
Description=vibmon IMU gateway
Documentation=https://github.com/wight-u/vibmon-gateway
# Start after the kernel module has registered /dev/vibmon0 and udev settled
After=dev-vibmon0.device systemd-udev-settle.service
BindsTo=dev-vibmon0.device

[Service]
Type=notify
ExecStart=/usr/bin/vibmon-gateway /dev/vibmon0 /dev/ttyAMA0
Restart=on-failure
RestartSec=5s

# Watchdog: gateway must call sd_notify(0, "WATCHDOG=1") at least every 15 s
WatchdogSec=30s
NotifyAccess=main

# Runtime data (SQLite DB, FFTW wisdom)
StateDirectory=vibmon
StateDirectoryMode=0750

User=vibmon
Group=vibmon
SupplementaryGroups=dialout

# Harden the service
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/vibmon
ProtectHome=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
